# ops/linux/systemd/sg-bot.service
[Unit]
Description=School Gate Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/school-gate/current
EnvironmentFile=/opt/school-gate/current/.env
ExecStart=/usr/bin/pnpm --filter @school-gate/bot start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
